# xiaoxi-mesh/xiaoxi-mesh.service
[Unit]
Description=小希-Mesh 消息中转服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xiaoxi-mesh
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/xiaoxi-mesh/logs/mesh.log
StandardError=append:/opt/xiaoxi-mesh/logs/mesh.log

[Install]
WantedBy=multi-user.target
